# build/files/supervisord.conf
[supervisord]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:main]
command=/usr/local/app/avalanchego --network-peer-list-num-validator-ips=10 --http-host=0.0.0.0 --public-ip-resolution-service=ifconfigCo --chain-config-dir=/config --http-allowed-hosts="*"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true



